# deploy/com.investment.inbox.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.investment.inbox</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/franciscojaviernavarro/.pyenv/versions/3.12.2/bin/python3</string>
        <string>-u</string>
        <string>/Users/franciscojaviernavarro/Documents/Claude_Projects/investment-agents/tools/process_inbox.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/franciscojaviernavarro/Documents/Claude_Projects/investment-agents</string>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/franciscojaviernavarro/.local/bin:/Users/franciscojaviernavarro/.nvm/versions/node/v20.11.0/bin:/Users/franciscojaviernavarro/.pyenv/versions/3.12.2/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/investment_inbox.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/investment_inbox.log</string>
</dict>
</plist>
